# backend/deploy/specto.service
[Unit]
Description=Specto Web Application
After=network.target postgresql.service

[Service]
Type=simple
User=specto
Group=specto
WorkingDirectory=/opt/specto
ExecStart=/opt/specto/specto server
Restart=always
RestartSec=5
Environment=SPECTO_HOST=0.0.0.0
Environment=SPECTO_PORT=8080
Environment=SPECTO_LOG_LEVEL=info
Environment=SPECTO_DB_PROVIDER=postgres
EnvironmentFile=-/etc/specto/specto.env
LimitNOFILE=65535
StandardOutput=journal
StandardError=journal
SyslogIdentifier=specto

[Install]
WantedBy=multi-user.target
